# roles/base/templates/etc/systemd/system/wpa_supplicant@wlan.service.d/customdependency.conf
#{{ansible_managed}}
[Unit]
After=systemd-nerwork.service
